# 2014:2015_seasons/file.cron
0 0 * * * python3 getData.py

0 30 * * * python3 predict.py 